# Ansible/roles/ns_server/templates/unbound-cache-load.service.j2
# Unbound cache load service (runs on boot)
# Сервис загрузки кэша Unbound (запускается при загрузке)
# {{ ansible_managed }}

[Unit]
Description=Load Unbound DNS cache from file
After=unbound.service
Requires=unbound.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'test -f {{ unbound_cache_dump_file }} && /usr/sbin/unbound-control load_cache < {{ unbound_cache_dump_file }} || true'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

